# systemd/google2snipeit.timer.template
[Unit]
Description=Google2Snipe-IT Sync Timer
Requires=google2snipeit.service

[Timer]
# Run 1 minute after boot if service was missed
OnBootSec=1min
# Run daily at 2:00 AM
OnCalendar=*-*-* 02:00:00
# Run immediately if timer was missed (e.g., system was off)
Persistent=true
# Unit to run
Unit=google2snipeit.service

[Install]
WantedBy=timers.target
