# systemd/google2snipeit.service.template
[Unit]
Description=Google2Snipe-IT ChromeOS Device Sync Service
After=network.target

[Service]
Type=oneshot
User=%SERVICE_USER%
WorkingDirectory=%WORKING_DIR%
EnvironmentFile=%WORKING_DIR%/.env
ExecStart=%VENV_BIN%/python %SCRIPT_DIR%/snipe-IT.py
StandardOutput=journal
StandardError=journal
